# crontab
*/5 * * * * /bin/sh /app/scripts/cronjobs/execute-scheduled-analysis.sh
*/10 * * * * /bin/sh /app/scripts/cronjobs/execute-experiment-poll.sh